# paccache.timer
[Unit]
Description=Clean-up old pacman pkg
[Timer]
OnCalendar=weekly
Persistent=true
[Install]
WantedBy=multi-user.target
